# mqtt-s7-connector/rootfs/etc/s6-overlay/s6-rc.d/mqtt-s7-connector/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: mqtt-s7-connector
# Runs mqtt-s7-connector
# ==============================================================================

bashio::log.info "Starting: mqtt-s7-connector"

exec /usr/bin/mqtt-s7-connector.sh
