# ha-mbmd/rootfs/etc/s6-overlay/s6-rc.d/init/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: ModBus Measurement Daemon
# Copies the example config files
# ==============================================================================
bashio::log.info "Copying example config files..."

cp -n /root/mbmd/mbmd.dist.yaml /config/
